# systemd/knock-monitor.service
[Unit]
Description=Knock-Knock SSH Honeypot and Monitor
After=network.target redis-server.service

[Service]
User=root
WorkingDirectory=/root/knock-knock
ExecStart=/root/knock-knock/.venv/bin/python -u monitor.py
# Add --save-knocks to store individual knock records in SQLite
# Add --max-knocks N to auto-add IPs to blocklist.txt after N total hits (disabled by default)
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
