# systemd/knock-monitor.service
[Unit]
Description=Knock-Knock SSH Honeypot and Monitor
After=network.target redis-server.service

[Service]
User=root
WorkingDirectory=/root/knock-knock
ExecStart=/root/knock-knock/.venv/bin/python -u monitor.py
# Add --save-knocks to store individual knock records in SQLite
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
