# systemd/knock-monitor.service
[Unit]
Description=Knock-Knock SSH Log Monitor
After=network.target redis-server.service

[Service]
User=root
WorkingDirectory=/root/knock-knock
# Point to the python binary inside your uv virtual environment
# Add --save-knocks to save full knock history to SQLite (useful for
# research queries, but the knocks table will grow to ~600MB+/year)
ExecStart=/root/knock-knock/.venv/bin/python monitor.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
