# systemd/knock-web.service
[Unit]
Description=Knock-Knock FastAPI Web Server
After=network.target redis-server.service knock-monitor.service

[Service]
User=root
WorkingDirectory=/root/knock-knock
# Plain HTTP (default)
ExecStart=/root/knock-knock/.venv/bin/python3 -m uvicorn main:app \
    --host 0.0.0.0 \
    --port 80 \
    --proxy-headers \
    --forwarded-allow-ips='*' \
    --workers 2 \
    --timeout-keep-alive 30 \
    --log-level warning
# HTTPS — use this ExecStart instead (remove the one above):
# ExecStart=/root/knock-knock/.venv/bin/python3 -m uvicorn main:app \
#     --host 0.0.0.0 \
#     --port 443 \
#     --ssl-keyfile /root/knock-knock/certs/key.pem \
#     --ssl-certfile /root/knock-knock/certs/cert.pem \
#     --proxy-headers \
#     --forwarded-allow-ips='*' \
#     --workers 2 \
#     --timeout-keep-alive 30 \
#     --log-level warning
# Environment=LOG_VISITORS=true  # Uncomment to log dashboard visitors to visitors.db
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
