# systemd/knock-monitor.service
[Unit]
Description=Knock-Knock Honeypot and Monitor
After=network.target redis-server.service

[Service]
User=root
WorkingDirectory=/root/knock-knock
EnvironmentFile=-/root/knock-knock/.env
ExecStart=/root/knock-knock/.venv/bin/python -u monitor.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
