# systemd/knock-web.service
[Unit]
Description=Knock-Knock FastAPI Web Server
After=network.target redis-server.service knock-monitor.service

[Service]
User=root
WorkingDirectory=/root/knock-knock
EnvironmentFile=-/root/knock-knock/.env
ExecStart=/root/knock-knock/.venv/bin/python main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
